# systemd/ploutosledger.service
[Unit]
Description=PloutosLedger API Server
After=network.target postgresql.service
Wants=network.target

[Service]
Type=forking
User=ploutos
Group=ploutos
WorkingDirectory=/opt/ploutosledger
Environment=NODE_ENV=production
EnvironmentFile=/opt/ploutosledger/.env
ExecStart=/usr/bin/pm2 resurrect
ExecReload=/usr/bin/pm2 reload all
ExecStop=/usr/bin/pm2 kill
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=ploutosledger

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ploutosledger/prisma/data /opt/ploutosledger/logs

[Install]
WantedBy=multi-user.target

